# noter/package/noter_0.1-1_amd64/DEBIAN/postinst
#!/bin/sh

mkdir /tmp/noter
chmod 777 -R /tmp/noter
systemctl daemon-reload
systemctl start noterd
systemctl enable noterd
