# noter-srv/package/noter-srv_0.1-1_amd64/DEBIAN/postinst
#!/bin/sh

mkdir /tmp/noter_srv
systemctl daemon-reload
systemctl start noter-srv
systemctl enable noter-srv
