# noter-srv/package/noter-srv_0.1-1_amd64/DEBIAN/prerm
#!/bin/sh

systemctl stop noter-srv
systemctl disable noter-srv
systemctl daemon-reload
systemctl reset-failed
rm -rf /tmp/noter_srv
